# deployment/systemd/system/mosmix-processor-s.timer
[Unit]
Description=Run mosmix-processor-l.service every hour

[Timer]
OnCalendar=*:29

[Install]
WantedBy=multi-user.target
