# deployment/mosmix-processor.timer
[Unit]
Description=Run mosmix-processor.service every 30 minutes

[Timer]
OnCalendar=*:0/30

[Install]
WantedBy=multi-user.target
